# etc/systemd/growatt.service
[Unit]
Description=Growatt Modbus Proxy & Telemetry Server
Documentation=https://github.com/remcovanmook/growatt
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=growatt
Group=growatt
WorkingDirectory=/opt/growatt
EnvironmentFile=/etc/default/growatt
ExecStart=/opt/growatt/.venv/bin/python /opt/growatt/growatt_server.py
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
