# docs/growatt-proxy.service
[Unit]
Description=Growatt Modbus Proxy and Telemetry Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=growatt
Group=growatt
WorkingDirectory=/opt/growatt-proxy
Environment="PATH=/opt/growatt-proxy/.venv/bin:/usr/local/bin:/usr/bin:/bin"

# Replace <INVERTER_IP> with your datalogger IP
ExecStart=/opt/growatt-proxy/.venv/bin/python growatt_server.py --device-ip <INVERTER_IP> --db /var/lib/growatt/growatt.db

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
